# unistatus.service
[Unit]
Description=Starts the unistatus program

[Serivce]
Type=oneshot
ExecStart=/usr/local/bin/unistatus

[Install]
WantedBy=multi-user.target